# roles/loadbalancer/templates/envoy.service.j2
[Unit]
Description=Envoy Proxy 

[Service]
TimeoutStartSec=0
Restart=always
ExecStart=/usr/bin/envoy \
  -c /etc/envoy/envoy.yaml

[Install]
WantedBy=multi-user.target
